# vestigia/ops/backup_schedule.cron
0 * * * * /home/arksher/ml-redteam/vestigia/backup.sh >> /var/log/vestigia_backup.log 2>&1
